# templates/fluent-bit.service.j2
[Unit]
Description=Fluentbit
Requires=network.target
After=network.target

[Service]
Type=simple
LimitNOFILE=65536
SupplementaryGroups=adm
UMask=0027
ExecStart=/usr/local/bin/fluent-bit -c /etc/fluent-bit/fluent-bit.conf
Restart=always

[Install]
WantedBy=multi-user.target
